# packages/embylogic/depends/emby4/system.d/emby4.service
[Unit]
Description=Emby 4 - streaming server
Documentation=http://emby.media
Wants=kodi.service
After=kodi.service
Conflicts=service.emby.service

[Service]
ExecStart=/bin/sh /usr/bin/emby4.start
Restart=always

[Install]
WantedBy=embylogic.service
